# rootfs/etc/s6-overlay/s6-rc.d/svc-msmtpd/run
#!/command/execlineb -P
with-contenv
msmtpd --port=25 --interface=127.0.0.1